# home/conf/bot_conf.conf
[program:bot]
command=/bot_integrations/venv/bin/python /bot_integrations/home/bot_father.py
environment=PATH=/bot_integrations/venv/bin
stdout_logfile=/bot_integrations/error.log
